# deploy/systemd/meteo.service
[Unit]
Description=Meteo SYNOP dashboard
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/meteo
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env npm start
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
